# botkadr.service
[Unit]
Description=Kadr
After=network.target

[Service]
EnvironmentFile=/etc/environment
ExecStart=/home/ubuntu/BotDavomad/venv/bin/python main.py
ExecReload=/home/ubuntu/BotDavomad/venv/bin/python main.py
WorkingDirectory=/home/ubuntu/BotDavomad/
KillMode=process
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target